# systemd/.config/systemd/user/set-bg.timer
[Unit]
Description=Change background every second

[Timer]
OnCalendar=*:0:0
Unit=set-bg.service

[Install]
WantedBy=default.target